# ansible/roles/kubernetes-common/templates/etc/systemd/system/kubelet.service.d/11-fix-dns.conf
[Service]
Environment="KUBELET_DNS_ARGS=--cluster-dns={{ ( kubernetes_common.service_cidr | ipaddr('subnet')) | ipaddr('10') | ipaddr('address') }} --cluster-domain=cluster.local"